# Postfix Stats/etc/cron.d/generate_pflogsumm
* * * * * root /etc/zabbix/zabbix_agent2.d/scripts/postfix_stats/generate_pflogsumm.sh